# mpbox_buildroot/fs/rootfs_customization/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system

# now run any rc scripts
::sysinit:/etc/init.d/rcS

# Put a getty or shell on the console
tty1::respawn:/bin/ash -l
#tty1::respawn:/sbin/getty 38400 tty1
#tty2::respawn:/sbin/getty 38400 tty2 

#Put a getty on the serial port
ttyS0::respawn:/sbin/getty -L ttyS0 115200 # GENERIC_SERIAL

# Start all of the mpservice stuff
null::respawn:/mpservice/bin/boa -d -c /mpservice/misc
null::respawn:/mpservice/bin/mDNSResponderPosix -f /var/tmp/mdns.conf
null::respawn:/mpservice/bin/mps-librarydbd > /mpservice/misc/mps-librarydbd.err 2>&1
null::respawn:/mpservice/bin/mps-server

# Start a telnetd
::respawn:/usr/sbin/telnetd -F

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/etc/init.d/rcK
null::shutdown:/bin/umount -r /disk1


